# .config/systemd/user/battery.timer
[Unit]
Description=Run battery check script every 60 seconds

[Timer]
OnBootSec=1min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
