# .config/systemd/user/battery.service
[Unit]
Description=Execute battery check script

[Service]
ExecStart=/usr/local/bin/battery
